# deployment/systemd/systemd-maintenance.service
# Systemd service pentru maintenance
# Copiază în /etc/systemd/system/stremio-maintenance.service

[Unit]
Description=Stremio addon maintenance
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/stremio-addon/maintenance.sh
User=root
StandardOutput=journal
StandardError=journal

# Timeout pentru a preveni blocarea
TimeoutStartSec=300

# Nu reporni automat - doar timer-ul decide când rulează
Restart=no